# wintria-env/wintria/wintria.cron
8,38,58 * * * * ~/webapps/windjango/apache2/bin/start
*/1 * * * * /bin/sh ~/webapps/windjango/Wintria/scripts/unpack.sh >> ~/webapps/windjango/Wintria/logs/unpack.log 2>&1
#*/25 * * * * ~/bin/indexer Article_article --rotate --config ~/webapps/windjango/Wintria/misc/sphinx.conf >> ~/webapps/windjango/Wintria/logs/searchd_log.txt
*/55 * * * * /home/wintrialucas/webapps/windjango/Wintria/scripts/safe_index.sh >> ~/webapps/windjango/Wintria/logs/searchd_log.txt
*/25 * * * * /bin/sh ~/webapps/windjango/Wintria/scripts/check_sphinx.sh >> ~/webapps/windjango/Wintria/logs/searchd_log.txt

0 */4 * * * /bin/sh ~/webapps/windjango/Wintria/scripts/push_trending.sh >> ~/webapps/windjango/Wintria/logs/push_trending.log 2>&1

@reboot !/bin/indexer Article_article --config ~/webapps/windjango/Wintria/misc/sphinx/sphinx.conf
@reboot ~/bin/searchd --config ~/webapps/windjango/Wintria/misc/sphinx/sphinx.conf